# scripts/vm/out/etc/containers/systemd/beszel-agent/beszel-agent.container
[Unit]
Description=beszel-agent container
After=local-fs.target network.target tailscaled.service
Wants=tailscaled.service
Requires=
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Restart=on-failure
RestartSec=5s



[Container]
Image=docker.io/henrygd/beszel-agent:latest
ContainerName=beszel-agent
Environment=PORT=45876
Environment=KEY=big key
Volume=/var/run/podman/podman.sock:/var/run/docker.sock:ro
SecurityLabelDisable=true
Network=host
Label=io.containers.autoupdate=registry

[Install]
# Start by default on boot
WantedBy=multi-user.target default.target
