# virter/out/testrepo_instance/hello@bar/quadlets/hello@bar.container
[Unit]
Description=Hello Service
Wants=network-online.target
After=network-online.target

[Container]
Image=docker.io/busybox:latest
Exec=/bin/sh -c "trap 'exit 0' INT TERM; while true; do echo Hello World; sleep 1; done"
Volume=hello.volume:/hellobar

[Install]
WantedBy=multi-user.target
