# virter/out/testrepo5/hello/quadlets/hello.container
[Unit]
Description=Hello Service
Wants=network-online.target
After=network-online.target

[Container]
ContainerName=busybox1
Image=docker.io/busybox:stable
Exec=/bin/sh -c "trap 'exit 0' INT TERM; while true; do echo Hello World; sleep 1; done"
Network=hello.network
Volume=hello.volume:/hello
Label=io.containers.autoupdate=registry

[Install]
WantedBy=multi-user.target
