# install/materia-update.timer
[Unit]
Description=Materia daily update

[Timer]
OnCalendar=*-*-* 05:00:00 America/New_York
Unit=materia@update.service

[Install]
WantedBy=timers.target
