# virter/out/testrepo_all/hello/data/hello_world.service
[Unit]
Description=Hello world script

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hello.sh
